# devops/docker/mongodb/supervisord-app.conf
[program:mongodb]
command=/usr/bin/mongod --dbpath /home/app --port 27017 --smallfiles
numprocs=1
autostart=true
autorestart=true
