# devops/docker/nodejs/supervisord-app.conf
[program:app]
command=/start.sh
numprocs=1
autostart=true
autorestart=true
stderr_logfile=/home/app/err.log
stdout_logfile=/home/app/out.log
